# deploy/resume-analyzer.service
[Unit]
Description=Resume Analyzer API
After=network.target

[Service]
Type=exec
User=hiros
Group=hiros
WorkingDirectory=/opt/resume-analyzer
EnvironmentFile=/opt/resume-analyzer/.env
ExecStart=/opt/resume-analyzer/backend/venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000 --workers 6
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/resume-analyzer/talent_pool.db /opt/resume-analyzer/talent-pool /opt/resume-analyzer/backend/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
